# api/corillo-api.service
[Unit]
Description=Corillo API — perfiles públicos y stream key management
After=network.target

[Service]
Type=simple
User=corillo-adm
WorkingDirectory=/home/corillo-adm/corillo-api
EnvironmentFile=/home/corillo-adm/corillo-api/.env
ExecStart=/home/corillo-adm/corillo-api/venv/bin/uvicorn server:app --host 127.0.0.1 --port 3004
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
